# .config/autostart/ctrlcaps.desktop
[Desktop Entry]
Name=ctrlcaps
Comment=Remap Capslock to Control
Type=Application
Exec=/usr/bin/setxkbmap -model evdev -option ctrl:nocaps
