# .config/autostart/smartcaps.desktop
[Desktop Entry]
Name=smartcaps
Comment=Smartcaps
Type=Application
Exec=smartcaps.xorg
